# crontab/CrontabeEntry
0 */4 * * * /home/notyours/Storage/dev/RealtorCAify/crontab/CrontabProcess.sh
